# xorg/.config/autostart/home.desktop
[Desktop Entry]
Type=Application
Name=Default File Manager
Exec=sh -c 'open "$HOME"'
