# xorg/.config/autostart/homepage.desktop
[Desktop Entry]
Type=Application
Name=Default Web Browser
Exec=sh -c 'open "$WWW_HOME"'
